# ops/systemd/sod-prayer-bot.service
[Unit]
Description=SOD Daily Prayer Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sod-emuna-agency
ExecStart=/opt/sod-emuna-agency/venv/bin/python sod/bots/prayer_bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
