# ops/systemd.service.example
[Unit]
Description=Emuna Agency Telegram sender (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/path/to/sod-emuna-agency
ExecStart=/usr/bin/python -m src.agency %i

[Install]
WantedBy=multi-user.target

# Timer example (save as /etc/systemd/system/agency@morning.timer etc.)
[Unit]
Description=Run Emuna Agency morning slot

[Timer]
OnCalendar=Mon..Fri 09:00
Unit=agency@morning.service

[Install]
WantedBy=timers.target
